# mock-server/config/10-xrootd-standalone.conf
[program:xrootd-standalone]
command=xrootd -c /etc/xrootd/xrootd-standalone.cfg -k fifo -s /tmp/xrootd/xrootd.pid
user=xrootd
autorestart=true
startsecs=10
